# deploy/kite-api.service
[Unit]
Description=Kite Trading API (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kite
ExecStart=/home/ubuntu/kite/venv/bin/uvicorn kite_app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
